# crontab.cron
# linux调度任务

* * * * *  php /data/www/blog/artisan schedule:run
